# aiglasses/systemd/aiglasses.target
[Unit]
Description=AI Glasses Platform Services
After=network-online.target bluetooth.target
Wants=network-online.target

[Install]
WantedBy=multi-user.target


